# agent/macos/com.reco3.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.reco3.agent</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/opt/reco3-agent/common/reco_agent.py</string>
        <string>/opt/reco3-agent/common/config.yaml</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/opt/reco3-agent</string>

    <key>StandardOutPath</key>
    <string>/var/log/reco3-agent/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/reco3-agent/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
